# control-plane/supervisord/hermes-control-plane-monitor.conf
[program:hermes-control-plane-monitor]
command=/workspace/hermes-agent-plugin/control-plane/supervisord/start-monitor.sh
directory=/workspace/hermes-agent-plugin
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stopwaitsecs=20
user=root
stdout_logfile=/opt/data/hermes/control-plane/monitor/supervisord.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
environment=HERMES_HOME="/opt/data/hermes",HERMES_CONTROL_PLANE_BOT_TOKEN_ENV="TELEGRAM_BOT_TOKEN",HERMES_CONTROL_PLANE_NOTIFY_CHAT_ENV="TELEGRAM_HOME_CHANNEL",HERMES_MONITOR_INTERVAL_SECONDS="10"
