# control-plane/supervisord/hermes-hersocial-marketing-automation.conf
[program:hermes-hersocial-marketing-automation]
command=/workspace/hermes-agent-plugin/control-plane/hersocial-marketing-automation/start-hersocial-marketing-automation.sh
directory=/workspace/hermes-agent-plugin
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stopwaitsecs=20
user=root
stdout_logfile=/opt/data/hermes/control-plane/hersocial-marketing-automation/supervisord.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
environment=HERSOCIAL_MARKETING_POLL_SECONDS="30"
