# control-plane/supervisord/hermes-hersocial-auto-post.conf
[program:hermes-hersocial-auto-post]
command=/workspace/hermes-agent-plugin/control-plane/hersocial-auto-post/start-hersocial-auto-post.sh
directory=/workspace/hermes-agent-plugin
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stopwaitsecs=20
user=root
stdout_logfile=/opt/data/hermes/control-plane/hersocial-auto-post/supervisord.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
environment=HERSOCIAL_AUTO_POST_ENABLED="false",HERSOCIAL_ATTENDED_REMINDERS_ENABLED="true",HERSOCIAL_AUTO_POST_POLL_SECONDS="30",HERSOCIAL_MARKETING_CHECK_COMMAND="sdtk-marketing",HERSOCIAL_MARKETING_CHECK_TIMEOUT_SECONDS="15"
